# epic.timer
[Unit]
Description=Run Epic games claimer weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
