# systemd/epic.service
[Unit]
Description=Run epic games claimer docker image

[Service]
Type=oneshot
ExecStart=docker start epic-claim
